# settings/celery.ini
[program:cron]
process_name=cron_service
directory=/path/to/common_module
command=celery worker -B -A celery_project.app --loglevel=info --autoscale=8,4
numprocs=1
user=work
startsecs=1
autorestart=true
autostart=true
startretries=5
redirect_stderr=true
stdout_logfile=/var/logs/celery/stdout.log
stdout_logfile_maxbytes=200MB
stdout_logfile_backups=5
stderr_logfile=/var/logs/celery/stderr.log
stderr_logfile_maxbytes=200MB
stderr_logfile_backups=5